# rpm/tracker-store.service
[Unit]
Description=Tracker Store
After=xorg.target
Requires=dbus.socket xorg.target

[Service]
ExecStart=/usr/libexec/tracker-store

